# .DOTFILES_RESOURCES/system/monitor-hotplug/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/bin/systemctl start monitorhotplug.service"
