# .DOTFILES_RESOURCES/system/monitor-hotplug/monitorhotplug.service
[Unit]
Description=Monitor hotplug

[Service]
Type=oneshot
RemainAfterExit=no
User=gorane
ExecStart=/usr/local/bin/monitor-hotplug.sh
