# .DOTFILES_RESOURCES/etc/systemd/system/i3lock.service.template
[Unit]
Description=i3lock on suspend
Before=sleep.target

[Service]
User=###USER###
Type=exec
Environment=DISPLAY=:0
ExecStart=/usr/local/bin/lock.sh

[Install]
WantedBy=sleep.target
